# geo-deb/src/deb/init/geosearch.conf
description "Geo Search Service"  
author "mygov.scot"

start on runlevel [2345]
stop on runlevel [!2345]

setuid geosearch
chdir /opt/geosearch
exec /opt/geosearch/run
