# geo-deb/src/deb/control/prerm
#!/bin/sh
systemctl stop    geosearch || true
systemctl disable geosearch || true
